# fintrack/rootfs/etc/services.d/fintrack/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting My Cash Supervisor..."

# ingress_port: 0 in config.yaml lässt den Supervisor den tatsächlichen Port
# dynamisch zuteilen (vermeidet Konflikte, falls 8099 anderswo belegt ist).
# Den zugewiesenen Port fragen wir hier ab, statt ihn fix anzunehmen.
PORT="$(bashio::addon.ingress_port)"
bashio::log.info "Ingress-Port: ${PORT}"

export PORT
export DB_PATH="/data/fintrack.db"

cd /app/server
exec node index.js
